# worker/launchd/com.vexecsearch.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Legacy daily plist — replaced by install_launchd.sh (scrape/rescore/enrich/imessage/email).
  Re-run: cd worker && ./scripts/install_launchd.sh
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vexecsearch.daily</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>PYTHON_PATH</string>
        <string>WORKER_ROOT/scripts/run_daily.py</string>
        <string>--scrape-only</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>WORKER_ROOT</string>
</dict>
</plist>
